# battery_warning/battery_warning.service.d/override.conf
# Documentation: `man 5 systemd.service`
# Documentation: `man 7 systemd.directives`

[Service]
User=hendrik
